# test683-daemontools/crontab
@reboot /usr/local/bin/svscan $HOME/service 
