# test680-gcs_fuse/provision/bucket.service.in
[Unit]
Description=Mount GCS bucket

[Service]
Type=simple
ExecStart=/usr/bin/gcsfuse --foreground {{storage_bucket}} %h/bucket
Restart=always

[Install]
WantedBy=default.target
